# mgmt/dnode/package/DEBIAN/prerm
systemctl stop unique-docker.service
systemctl disable /lib/systemd/system/unique-docker.service
